# projects/ROCKNIX/packages/rocknix/system.d/rocknix-memory-manager.service
[Unit]
Description=Rocknix Memory Manager
After=local-fs.target
DefaultDependencies=no
Conflicts=swap.service shutdown.target
Before=shutdown.target

[Service]
Environment=HOME=/storage
Type=oneshot
ExecStart=/usr/bin/rocknix-memory-manager --setup
ExecReload=/usr/bin/rocknix-memory-manager --reload
ExecStop=/usr/bin/rocknix-memory-manager --stop
RemainAfterExit=yes

[Install]
WantedBy=basic.target
