# projects/ROCKNIX/packages/tools/inputplumber/sources/usr/lib/systemd/system/inputplumber.service
[Unit]
Description=InputPlumber Service

[Service]
Type=dbus
BusName=org.shadowblip.InputPlumber
Environment=LOG_LEVEL=info
Environment=HIDE_DEVICES_FROM_ROOT=1
Environment=INSECURE_DISABLE_POLKIT=1
ExecStart=/usr/bin/inputplumber
ProtectSystem=full

[Install]
WantedBy=multi-user.target
