# projects/ROCKNIX/packages/hardware/quirks/system.d/volume-fixup.service
[Unit]
Description=Fix unset volumes on RK3566 devices 
After=rocknix.target
Before=essway.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/lib/autostart/quirks/platforms/RK3566/bin/volumealign

[Install]
WantedBy=rocknix.target
